# infra/systemd/vezir.service
[Unit]
Description=vezir — internal scribe service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Vezir is installed via `pip install --user -e .` against /usr/bin/python3
# (deps already present from meetscribe). Adjust ExecStart if you moved
# vezir into a venv.
ExecStart=%h/.local/bin/vezir serve
WorkingDirectory=%h/models/vezir
Environment=VEZIR_DATA=%h/vezir-data
Environment=VEZIR_HOST=0.0.0.0
Environment=VEZIR_PORT=8000
# PATH includes ~/.local/bin (for `meet`), and the standard system
# locations for ffmpeg, git, etc.
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Secrets (HF_TOKEN, OPENROUTER_API_KEY, MEETSCRIBE_SUMMARY_BACKEND).
# File is mode 0600. Leading '-' means it is optional.
EnvironmentFile=-%h/.config/environment.d/vezir.conf
UMask=0077
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
